# deploy/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/app/logs/supervisord.log
pidfile=/var/run/supervisord.pid

[program:gunicorn]
command=/opt/venv/bin/gunicorn --bind 127.0.0.1:8000 --workers 4 --timeout 600 --access-logfile /app/logs/gunicorn-access.log --error-logfile /app/logs/gunicorn-error.log tabgraphsyn_site.wsgi:application
directory=/app
user=root
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/app/logs/gunicorn.log
environment=PATH="/opt/venv/bin"

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/app/logs/nginx.log
